# systemd/netmon-server.service
[Unit]
Description=netmon echo server (TCP+UDP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=netmon
Group=netmon
EnvironmentFile=/etc/netmon/netmon.conf
ExecStart=/usr/local/bin/netmon-server --host 0.0.0.0 --port ${REMOTE_PORT}
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
