# systemd/netmon-microcuts.service
[Unit]
Description=netmon microcuts monitor (UDP loss detection)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=netmon
Group=netmon
EnvironmentFile=/etc/netmon/netmon.conf
ExecStart=/usr/local/bin/netmon-microcuts \
    --host ${REMOTE_HOST} \
    --port ${REMOTE_PORT} \
    --rate ${MICROCUTS_RATE} \
    --window ${MICROCUTS_WINDOW} \
    --burst-threshold ${MICROCUTS_BURST_THRESHOLD} \
    --log-dir /var/log/netmon/microcuts
Restart=always
RestartSec=10

# Prioridad alta para precisión de timing
Nice=-5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/netmon
ProtectHome=true

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
